# devops/roles/supervisor/templates/app.conf.j2
[program:{{ app_name }}]
command = {{ virtualenv_path }}/bin/uwsgi {{ virtualenv_path }}/uwsgi.ini
user = {{ uwsi_user }}
stdout_logfile = {{ application_log_file }}
redirect_stderr = true