# devops/roles/supervisor/templates/worker.conf.j2
[program:{{ app_name }}-worker]
command = {{ virtualenv_path }}/bin/rqworker {{ app_name }}
user = {{ gunicorn_user }}
directory = {{ project_path }}