# devops/roles/supervisor/templates/worker.conf.j2
[program:{{ application_name }}-worker]
command = {{ virtualenv_path }}/bin/rqworker {{ application_name }}
user = {{ app_user }}
directory = {{ project_path }}